# deploy/tencent/buddys.service
[Unit]
Description=Buddys Runtime API
After=network.target
# Reverse proxy target: 127.0.0.1:8787

[Service]
Type=simple
User=root
WorkingDirectory=/opt/buddys
ExecStart=/opt/buddys/deploy/tencent/run_with_env_compat.sh /etc/buddys/buddys.env /usr/bin/env PYTHONPATH=/opt/buddys/src /opt/buddys/.venv/bin/python -m uvicorn buddys_api.main:create_app --factory --host 127.0.0.1 --port 8787
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
